# .config/systemd/user/backup.timer
[Unit]
Description=Run backup each day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
