# .config/systemd/user/backup-server.service
[Unit]
Description=Server backup

[Service]
Type=oneshot
ExecStart=/home/seb/.bin/backup-server
