# .config/systemd/user/backup-server.timer
[Unit]
Description=Run backup-server each day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
