# .config/systemd/user/backup.service
[Unit]
Description=Notebook backup

[Service]
Type=oneshot
ExecStart=/home/blinry/.bin/backup
