# .config/systemd/user/doing.service
[Service]
Type=simple
ExecStart=%h/.bin/doing
Environment="DISPLAY=:0" "LANG=en_US.UTF-8"
