# .config/systemd/user/offlineimap.service
[Unit]
Description=Offlineimap sync

[Service]
Type=oneshot
ExecStart=offlineimap
