# .config/systemd/user/offlineimap.timer
[Unit]
Description=Run offlineimap each hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
