# .config/systemd/user/web-block.service
[Unit]
Description=Block distractions

[Service]
Type=oneshot
ExecStart=/home/blinry/.bin/web-block
